# ialarm-mqtt/rootfs/etc/services.d/ialarm-mqtt/finish
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Custom Add-on: ialarm-mqtt
# Take down the S6 supervision tree when ialarm-mqtt fails
# ==============================================================================
if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
  bashio::log.warning "ialarm-mqtt crashed, halting add-on"
  /run/s6/basedir/bin/halt
fi

bashio::log.info "ialarm-mqtt stopped, restarting..."
